# full_install/restic/restic-backup.timer
[Unit]
Description=Run Restic Offsite Backup Daily

[Timer]
# Schedule the backup for 2:00 AM daily
OnCalendar=*-*-* 02:00:00

# Run immediately if the system was down when the timer should have fired
Persistent=true 

[Install]
WantedBy=timers.target
